# conf/supervisor-collectlogs.conf
[program:collogs]
environment=LANG="en_US.UTF-8"
command=python /opt/servplus/ds/devtools/collogs/start.py

